# timers/osmen-secrets-audit.timer
[Unit]
Description=OsMEN-OC secrets audit (daily)
Documentation=https://github.com/dwill/OsMEN-OC

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=300

[Install]
WantedBy=timers.target
