# timers/osmen-chromadb-compact.service
# osmen-chromadb-compact.service
# One-shot service: compacts ChromaDB collections to reclaim space.
# Invoked by osmen-chromadb-compact.timer.

[Unit]
Description=OsMEN-OC — ChromaDB compaction
Documentation=https://github.com/dwilli15/OsMEN-OC

[Service]
Type=oneshot
ExecStart=/home/dwill/dev/OsMEN-OC/.venv/bin/python3 /home/dwill/dev/OsMEN-OC/scripts/chromadb-compact.py
TimeoutStartSec=300

[Install]
WantedBy=default.target
