# timers/osmen-secrets-audit.service
[Unit]
Description=OsMEN-OC secrets audit runner
Documentation=https://github.com/dwill/OsMEN-OC

[Service]
Type=oneshot
EnvironmentFile=%h/.config/osmen/env
ExecStart=%h/dev/.venv/bin/python -m core.secrets.cli audit
WorkingDirectory=%h/dev/OsMEN-OC
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
