# quadlets/tools/osmen-tools-bentopdf.container
# osmen-tools-bentopdf.container
# Rootless Podman Quadlet — BentoPDF self-hosted PDF toolkit.
# Privacy-first, client-side processing.

[Unit]
Description=OsMEN-OC Tools — BentoPDF
Documentation=https://github.com/alam00000/bentopdf

[Container]
Image=ghcr.io/alam00000/bentopdf:latest
ContainerName=osmen-tools-bentopdf

PublishPort=127.0.0.1:3020:8080

HealthCmd=wget -q --spider http://localhost:8080 || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=15s

SecurityLabelDisable=true
NoNewPrivileges=true

[Service]
Restart=always
RestartSec=5s
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=default.target
