# quadlets/dashboard/user-osmen-dashboard.slice
[Unit]
Description=OsMEN-OC Dashboard Services
Before=slices.target

[Slice]
MemoryMax=512M
CPUQuota=100%
