# quadlets/media/osmen-media-gluetun.container
# osmen-media-gluetun.container
# Rootless Podman Quadlet — gluetun VPN gateway (WireGuard).
# Runs inside download-stack pod, providing network namespace.
# Uses "custom" provider because gluetun doesn't natively support Privado.
# PF04: Verify DNS doesn't leak (cat /etc/resolv.conf inside container).
# PF05: Verify IPv6 disabled (curl -6 ifconfig.me must FAIL).
# PF06: VPN creds via Podman secrets, never in git.
# PF07: Must auto-start on reboot.

[Unit]
Description=OsMEN-OC Media — gluetun VPN Gateway
Documentation=https://github.com/dwilli15/OsMEN-OC

[Container]
Image=docker.io/qmcgaw/gluetun:v3.40.2
ContainerName=osmen-media-gluetun
Pod=download-stack.pod

# Custom WireGuard provider (Privado VPN Amsterdam).
Environment=VPN_SERVICE_PROVIDER=custom
Environment=VPN_TYPE=wireguard
Environment=VPN_ENDPOINT_IP=91.148.236.70
Environment=VPN_ENDPOINT_PORT=51820
Environment=WIREGUARD_PUBLIC_KEY=KgTUh3KLijVluDvNpzDCJJfrJ7EyLzYLmdHCksG4sRg=
Environment=DNS_ADDRESS=198.18.0.1

# Disable IPv6 (PF05).
Environment=DOT_IPV6=off

# VPN credentials via Podman secrets.
Secret=osmen-vpn-wireguard-private-key,type=env,target=WIREGUARD_PRIVATE_KEY
Secret=osmen-vpn-wireguard-addresses,type=env,target=WIREGUARD_ADDRESSES

# NET_ADMIN required for VPN tunnel creation.
AddCapability=NET_ADMIN
AddDevice=/dev/net/tun

HealthCmd=wget -q -O /dev/null http://127.0.0.1:9999
HealthInterval=60s
HealthTimeout=15s
HealthRetries=5
HealthStartPeriod=30s

SecurityLabelDisable=true
NoNewPrivileges=true
# ReadOnly disabled — gluetun writes to /etc/passwd and /tmp at startup.
Tmpfs=/tmp
Sysctl=net.ipv6.conf.all.disable_ipv6=1

[Service]
Restart=always
RestartSec=10s
Slice=user-osmen-media.slice
MemoryMax=256M
CPUQuota=50%

[Install]
WantedBy=default.target
