# quadlets/core/osmen-core-postgres.container
# osmen-core-postgres.container
# Rootless Podman Quadlet — PostgreSQL 17 with pgvector extension.
# Exposure: 127.0.0.1:5432 only (localhost-only initial policy).
# Secrets managed via Podman secrets (SOPS-decrypted at deploy time).

[Unit]
Description=OsMEN-OC Core — PostgreSQL 17 + pgvector
<<<<<<< HEAD
Documentation=https://github.com/dwilli15/OsMEN-OC
=======
Documentation=https://github.com/OsMEN-OC/OsMEN-OC
>>>>>>> origin/main
After=osmen-core-network.service
Requires=osmen-core-network.service

[Container]
# pgvector/pgvector bundles pg17 + the pgvector extension in one image.
Image=docker.io/pgvector/pgvector:pg17
ContainerName=osmen-core-postgres

# Internal-only: bind to loopback so no LAN exposure in initial phase.
PublishPort=127.0.0.1:5432:5432

Network=osmen-core.network

# Persistent data volume — created automatically by Podman if absent.
Volume=osmen-postgres-data.volume:/var/lib/postgresql/data

# Credentials via Podman secrets (never plaintext env vars).
Secret=osmen-postgres-password,type=env,target=POSTGRES_PASSWORD
Secret=osmen-postgres-user,type=env,target=POSTGRES_USER
Secret=osmen-postgres-db,type=env,target=POSTGRES_DB

# Enable WAL-level logical replication for future CDC use.
Exec=postgres -c wal_level=logical -c max_connections=200

<<<<<<< HEAD
# Use fixed identifiers here so Quadlet can emit a direct health command cleanly.
HealthCmd=pg_isready -U osmen -d osmen
=======
HealthCmd=/bin/sh -c "pg_isready -U $POSTGRES_USER -d $POSTGRES_DB"
>>>>>>> origin/main
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5
HealthStartPeriod=30s

ReadOnly=true
SecurityLabelDisable=true
NoNewPrivileges=true

[Service]
Restart=always
RestartSec=5s
Slice=user-osmen-core.slice
MemoryMax=1G
CPUQuota=200%

[Install]
WantedBy=default.target
