# quadlets/monitoring/osmen-monitoring-uptimekuma.container
# osmen-monitoring-uptimekuma.container
# Rootless Podman Quadlet — Uptime Kuma (service health monitoring).
# Exposure: 127.0.0.1:3001 only (behind Caddy: uptimekuma.osmen.local).
# Alert channels: Telegram (primary), Discord (secondary).

[Unit]
Description=OsMEN-OC Monitoring — Uptime Kuma
Documentation=https://github.com/dwilli15/OsMEN-OC
After=osmen-core-network.service
Requires=osmen-core-network.service

[Container]
Image=docker.io/louislam/uptime-kuma:1.23.16
ContainerName=osmen-monitoring-uptimekuma

PublishPort=127.0.0.1:3001:3001

Network=osmen-core.network

Volume=osmen-uptimekuma-data.volume:/app/data

HealthCmd=curl -sf http://127.0.0.1:3001
HealthInterval=30s
HealthTimeout=10s
HealthRetries=5
HealthStartPeriod=15s

SecurityLabelDisable=true
NoNewPrivileges=true
ReadOnly=true

[Service]
Restart=always
RestartSec=5s
Slice=user-osmen-background.slice
MemoryMax=512M
CPUQuota=100%

[Install]
WantedBy=default.target
