# quadlets/inference/osmen-inference-lmstudio.service
# osmen-inference-lmstudio.service
# Systemd user service — LM Studio headless (Vulkan inference on AMD Radeon 780M iGPU).
# NOT a container — native binary installed to ~/.lmstudio/bin/.
# Exposes OpenAI-compatible API on 127.0.0.1:1234.
# PF11: When FFXIV uses NVIDIA GPU, inference routes to AMD Vulkan via this.

[Unit]
Description=OsMEN-OC Inference — LM Studio Headless (Vulkan/AMD 780M)
Documentation=https://github.com/dwilli15/OsMEN-OC
# No network dependency — this is a native binary, not a container.

[Service]
Type=simple
ExecStart=%h/.lmstudio/bin/lms server start --port 1234
ExecStop=%h/.lmstudio/bin/lms server stop
Restart=on-failure
RestartSec=10s
Slice=user-osmen-inference.slice
MemoryMax=8G
CPUQuota=400%

# Vulkan environment for AMD iGPU.
Environment=VULKAN_SDK=/usr
Environment=VK_ICD_FILENAMES=/usr/share/vulkan/icd.d/radeon_icd.json

[Install]
WantedBy=default.target
